# config/crontab
* * * * * PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin; /bin/date >> /tmp/cronlog && /bin/mkdir /tmp/hoge
