# config/crontab
* * * * * bash -c "date >> /tmp/cronlog && mkdir /tmp/hoge"
